# 99-pendant-backlight.rules
# Pendant 키오스크: 백라이트 brightness 파일을 video 그룹에서 쓰기 가능하게.
# (이전엔 부팅마다 max_brightness 를 강제로 써 밝기를 최대 고정했으나 제거.
#  밝기는 앱이 settings.json 의 "screen_brightness" 값으로 설정한다.)
# odroid 사용자는 video 그룹이므로 sudo 없이 밝기 조절 가능.
# setup.sh 가 /etc/udev/rules.d/ 로 설치.
ACTION=="add", SUBSYSTEM=="backlight", KERNEL=="backlight", RUN+="/bin/chgrp video /sys/class/backlight/%k/brightness", RUN+="/bin/chmod g+w /sys/class/backlight/%k/brightness"
