# pendant.service
[Unit]
Description=Pendant HMI (M1S, weston kiosk)
After=multi-user.target systemd-user-sessions.service
Conflicts=getty@tty1.service

[Service]
Type=simple
User=odroid
Group=odroid
SupplementaryGroups=video render input dialout gpio netdev
WorkingDirectory=/home/odroid/hana-Pendant

Environment=HOME=/home/odroid
Environment=XDG_RUNTIME_DIR=/run/user/1000

; tty1 점유 + 콘솔 cleanup (getty 와 충돌)
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

# systemd-logind 세션 생성 → seat0 활성 → weston 이 DRM/input takeover 가능
# (이 옵션이 launch.sh stdout 을 tty 로 뺏아가므로 launch.sh 가 자체 파일로 redirect)
PAMName=login

ExecStart=/home/odroid/hana-Pendant/launch.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
